# ai_trading_bot.service
[Unit]
Description=AI Trading Bot Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/ai_trading_bot
Environment="PATH=/home/ubuntu/ai_trading_bot/venv/bin:/usr/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/ubuntu/ai_trading_bot/venv/bin/python ai_trading_bot.py
StandardOutput=append:/home/ubuntu/ai_trading_bot/bot.log
StandardError=append:/home/ubuntu/ai_trading_bot/bot.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
